# ansible/roles/health_service/templates/rewards-health.service.j2
[Unit]
Description=Rewards Health Service
Documentation=https://github.com/nealstreet/rewards-tier
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ health_service_app_user }}
Group={{ health_service_app_group }}
WorkingDirectory={{ health_service_app_dir }}

# Refresh secrets as root before starting the unprivileged service.
ExecStartPre=+{{ health_service_app_dir }}/fetch-secrets.sh

# Optional until the first successful secret fetch.
EnvironmentFile=-{{ health_service_app_dir }}/.env

ExecStart=/usr/bin/python3 {{ health_service_app_dir }}/health-service.py

Restart=always
RestartSec=10

NoNewPrivileges=true
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ health_service_service_name }}

[Install]
WantedBy=multi-user.target
